# kylin/ja_JP/airootfs/etc/profile.d/packey-reset.sh
# reset pacman's keyrings

rm -rf /etc/pacman.d/gnupg
pacman-key --init
pacman-key --populate archlinux
