# xfce/ja_JP/airootfs/etc/profile.d/localeset.csh
# set locale.
echo "LANG=$(cat /etc/locale.conf|grep LANG=|grep -oP '(?<==).+')"|tee /etc/locale.conf
